# adventureland-node-eu.service
[Unit]
Description=Adventure Land Node.js Game Server (EU Region)
After=network.target adventureland-backend.service
Requires=adventureland-backend.service

[Service]
Type=simple
User=youruser
WorkingDirectory=/opt/adventureland
ExecStart=/usr/bin/node /opt/adventureland/node/server.js EU I 2053
Restart=always
RestartSec=10
StandardOutput=append:/opt/adventureland/logs/node-eu.log
StandardError=append:/opt/adventureland/logs/node-eu-error.log

[Install]
WantedBy=multi-user.target
